# systemd/download-provider.socket
[Socket]
ListenStream=/tmp/.download-provider.sock
SocketMode=0777
PassCredentials=yes
Accept=false

[Install]
WantedBy=sockets.target
