# systemd/download-provider.socket
[Socket]
ListenStream=/tmp/download-provider
SocketMode=0777

[Install]
WantedBy=sockets.target
